# tools/arvbox/lib/arvbox/docker/service/banner/run
#!/bin/sh
# Copyright (C) The Arvados Authors. All rights reserved.
#
# SPDX-License-Identifier: AGPL-3.0

set -e

/usr/local/lib/arvbox/runsu.sh $0-service $1
sv stop banner
