# services/dockercleaner/arvados-docker-cleaner.service
# Copyright (C) The Arvados Authors. All rights reserved.
#
# SPDX-License-Identifier: AGPL-3.0

[Unit]
Description=Arvados Docker Image Cleaner
Documentation=https://doc.arvados.org/
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10s
RestartPreventExitStatus=2
ExecStart=/usr/bin/arvados-docker-cleaner

[Install]
WantedBy=multi-user.target
