# arvbox/lib/arvbox/docker/service/vm/run
#!/bin/bash
set -e

. /usr/local/lib/arvbox/common.sh

git config --system "credential.http://$localip:${services[arv-git-httpd]}/.username" none
git config --system "credential.http://$localip:${services[arv-git-httpd]}/.helper" '!cred(){ cat >/dev/null; if [ "$1" = get ]; then echo password=$ARVADOS_API_TOKEN; fi; };cred'

exec /usr/local/lib/arvbox/runsu.sh $0-service
