# tools/arvbox/lib/arvbox/docker/service/slurmd/run
#!/bin/bash

exec 2>&1
set -eux -o pipefail

exec /usr/local/lib/arvbox/runsu.sh /usr/sbin/slurmd -v -D
