# arvbox/lib/arvbox/docker/service/sdk/run
#!/bin/sh
set -e

/usr/local/lib/arvbox/runsu.sh $0-service
sv stop sdk
