# docker/slurm/supervisor.conf
[program:munge]
user=root
command=/etc/init.d/munge start

[program:slurm]
user=root
command=/etc/init.d/slurm-llnl start
