# docker/warehouse/supervisor.conf
[program:mysqld]
user=mysql
command=/usr/sbin/mysqld --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib/mysql/plugin --user=mysql --pid-file=/var/run/mysqld/mysqld.pid --socket=/var/run/mysqld/mysqld.sock --port=3306

[program:keepd]
user=root
command=/usr/src/warehouse-apps/libwarehouse-perl/server/keepd
