# docker/warehouse/supervisor.conf
[program:keepd]
user=root
command=/usr/src/warehouse-apps/libwarehouse-perl/server/keepd
